# .config/autostart/unclutter.desktop
[Desktop Entry]
Exec=unclutter
Icon=
Name=unclutter
Path=
Terminal=True
Type=Application
